# 11task/debpackage/debpackage/package/DEBIAN/cron.d
#
# Launches super.sh periodically
#
*/5 * * * * /usr/bin/scrlogrot
