# Scr_4/screen/screen/package/etc/systemd/system/SCREENPAD.timer
[Unit]
Description=time

[Timer]
OnBootSec=40

[Install]
WantedBy=default.target
