# 11task_correct/debpackage/debpackage/package/etc/systemd/system/ADVO.service
[Unit]
Description=Script

[Service]
ExecStart=/usr/bin/scriptlogs

[Install]
WantedBy=multi-user.target
