# 11task/debpackage/debpackage/package/usr/bin/VODA.service
[Unit]
Description=Script

[Service]
ExecStart=/usr/bin/scrlogrot

[Install]
WantedBy=multi-user.target
